# ansible/roles/nextcloud-server/templates/cron.d/nextcloud-cron.j2
*/15 * * * * root /usr/bin/docker exec -u www-data nextcloud-apache php /var/www/html/cron.php
